# include/rhel9/SOURCES/siptools-research.timer
[Timer]
OnActiveSec=0s
OnUnitActiveSec=60s

[Install]
WantedBy=timers.target
